# openvpn-v2ray/openvpn.ini
[program:openvpn]
command = /usr/sbin/openvpn --config /root/client.ovpn --auth-nocache
autostart = true
startsecs = 5
autorestart = true
user = root
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
